# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
pidfile=/tmp/supervisord.pid

[program:web]
command=node /app/server.js
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
environment=PORT="3000",HOSTNAME="0.0.0.0"

[program:worker]
command=python3 /app/apps/worker-py/main.py
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:scheduler]
command=python3 /app/apps/worker-py/scheduler.py
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
